# pi-player/services/pi-player.service
[Unit]
Description=Pi Player API Server
Documentation=https://github.com/your-org/pi-player
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-player
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/pi/pi-player
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/pi-player/pi_server.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-player

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/pi/pi-player
CapabilityBoundingSet=

# Resource limits
LimitNOFILE=1024
LimitNPROC=512

[Install]
WantedBy=multi-user.target