# services/media-player.service
[Unit]
Description=Pi Player Media Player Daemon
Documentation=https://github.com/rgsuhas/connect
After=graphical-session.target sound.target
Wants=graphical-session.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/connect
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/pi/connect
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=PULSE_SERVER=tcp:localhost:4713
ExecStart=/bin/bash -c 'if [ -f /home/pi/connect/venv/bin/python ]; then /home/pi/connect/venv/bin/python /home/pi/connect/media_player.py; else /usr/bin/python3 /home/pi/connect/media_player.py; fi'
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=media-player

# Security settings (less restrictive for media playback)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/pi/connect
ReadOnlyPaths=/usr/bin/vlc /usr/bin/cvlc /usr/bin/feh /usr/bin/ffmpeg

# Resource limits
LimitNOFILE=1024
LimitNPROC=512

# Allow access to audio/video devices
DeviceAllow=/dev/dri rw
DeviceAllow=/dev/snd rw
DeviceAllow=char-alsa rw

[Install]
WantedBy=graphical-session.target